# home/dot_config/systemd/user/transmission-daemon.service.tmpl
{{- /* Transmission daemon systemd service - VPN bound */ -}}
[Unit]
Description=Transmission BitTorrent Daemon (VPN-bound)
Documentation=man:transmission-daemon(1)
After=network-online.target
Wants=network-online.target
# VPN killswitch: transmission-vpn-bind script ensures VPN is connected before starting
# Manual stop required if VPN disconnects (user services can't BindsTo system services)

[Service]
Type=notify
# Update bind address before starting
ExecStartPre={{ .chezmoi.homeDir }}/.local/bin/transmission-vpn-bind
# Start daemon in foreground mode
ExecStart=/usr/bin/transmission-daemon --foreground --config-dir {{ .chezmoi.homeDir }}/.config/transmission-daemon
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID

# Restart on failure but not too quickly
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ .chezmoi.homeDir }}/Videos {{ .chezmoi.homeDir }}/.config/transmission-daemon {{ .chezmoi.homeDir }}/Downloads

# Resource limits
CPUQuota=50%
MemoryMax=500M

[Install]
WantedBy=default.target